# ybarojs-debian/DEBIAN/postinst
#!/bin/sh
set -e

# postinst: Post-installation tasks.
# On configure: install JetBrainsMono Nerd Font and tuigreet binary if not already present.
# Always: rebuild the system font cache.

FONT_DIR="/usr/local/share/fonts/nerd-fonts/JetBrainsMono"
FONT_ZIP_URL="https://github.com/ryanoasis/nerd-fonts/releases/download/v3.2.1/JetBrainsMono.zip"

case "$1" in
    configure)
        if fc-list | grep -qi "JetBrainsMono Nerd"; then
            echo "postinst: JetBrainsMono Nerd Font already installed, skipping download."
        else
            if command -v wget >/dev/null 2>&1; then
                echo "postinst: JetBrainsMono Nerd Font not found, installing..."
                mkdir -p "$FONT_DIR"

                # Use mktemp to avoid /tmp symlink attack (CVE class: TOCTOU race)
                FONT_ZIP_TMP="$(mktemp /tmp/JetBrainsMono.XXXXXX.zip)"

                echo "postinst: downloading JetBrainsMono.zip from GitHub Releases..."
                wget -q --tries=3 -O "$FONT_ZIP_TMP" "$FONT_ZIP_URL" || {
                    echo "postinst: WARNING: font download failed (network unavailable?). Waybar icon glyphs may not render correctly."
                    echo "postinst: To install manually later: dpkg-reconfigure ybarojs-debian"
                    rm -f "$FONT_ZIP_TMP"
                    FONT_ZIP_TMP=""
                    true
                }

                if [ -n "$FONT_ZIP_TMP" ] && [ -f "$FONT_ZIP_TMP" ]; then
                    echo "postinst: extracting .ttf files to $FONT_DIR..."
                    unzip -o -q "$FONT_ZIP_TMP" '*.ttf' -d "$FONT_DIR" || true

                    echo "postinst: cleaning up temporary zip..."
                    rm -f "$FONT_ZIP_TMP"

                    echo "postinst: rebuilding font cache for $FONT_DIR..."
                    fc-cache -f "$FONT_DIR" || true

                    echo "postinst: JetBrainsMono Nerd Font installed successfully."
                fi
            else
                echo "postinst: WARNING: wget not found. Cannot download JetBrainsMono Nerd Font."
                echo "postinst: Install wget and run: dpkg-reconfigure ybarojs-debian"
            fi
        fi

        echo "postinst: rebuilding system font cache..."
        fc-cache -f || true

        # Download and install impala WiFi management tool
        echo ""
        echo "postinst: ==============================================="
        echo "postinst: Installing impala WiFi management tool"
        echo "postinst: ==============================================="

        IMPALA_DIR="/usr/local/bin"
        IMPALA_BIN="${IMPALA_DIR}/impala"
        IMPALA_API_URL="https://api.github.com/repos/pythops/impala/releases/latest"

        if [ -f "$IMPALA_BIN" ]; then
            echo "postinst: impala already installed, skipping download."
        else
            if command -v wget >/dev/null 2>&1 || command -v curl >/dev/null 2>&1; then
                echo "postinst: Downloading latest impala release..."
                
                # Get download URL from GitHub API
                if command -v curl >/dev/null 2>&1; then
                    IMPALA_URL=$(curl -sL "$IMPALA_API_URL" | grep -o '"browser_download_url": "[^"]*x86_64-unknown-linux-musl"' | head -n 1 | cut -d'"' -f4)
                else
                    IMPALA_URL=$(wget -qO- "$IMPALA_API_URL" | grep -o '"browser_download_url": "[^"]*x86_64-unknown-linux-musl"' | head -n 1 | cut -d'"' -f4)
                fi

                if [ -n "$IMPALA_URL" ]; then
                    mkdir -p "$IMPALA_DIR"
                    
                    if command -v curl >/dev/null 2>&1; then
                        curl -sL "$IMPALA_URL" -o "$IMPALA_BIN" 2>/dev/null || true
                    else
                        wget -q "$IMPALA_URL" -O "$IMPALA_BIN" 2>/dev/null || true
                    fi

                    if [ -f "$IMPALA_BIN" ] && [ -s "$IMPALA_BIN" ]; then
                        chmod +x "$IMPALA_BIN"
                        echo "postinst: impala installed successfully to $IMPALA_BIN"
                        echo "postinst: Run 'impala' to manage WiFi networks (requires iwd)."
                    else
                        echo "postinst: WARNING: impala download failed. WiFi management will not be available."
                        echo "postinst: To install manually later, visit: https://github.com/pythops/impala"
                        rm -f "$IMPALA_BIN"
                    fi
                else
                    echo "postinst: WARNING: Could not find impala download URL."
                    echo "postinst: You can install it manually from: https://github.com/pythops/impala"
                fi
            else
                echo "postinst: WARNING: Neither curl nor wget available. Cannot download impala."
                echo "postinst: Install curl or wget, then reinstall this package."
            fi
        fi

        echo "postinst: ==============================================="
        echo ""

        # Download and install sysc-greet display manager greeter
        echo ""
        echo "postinst: ==============================================="
        echo "postinst: Installing sysc-greet display manager greeter"
        echo "postinst: ==============================================="

        SYSC_BIN="/usr/local/bin/sysc-greet"
        SYSC_SHARE_DIR="/usr/local/share/sysc-greet"
        SYSC_VERSION="1.1.6"
        SYSC_DEB_URL="https://github.com/Nomadcxx/sysc-greet/releases/download/v${SYSC_VERSION}/sysc-greet_${SYSC_VERSION}_amd64.deb"

        if [ -f "$SYSC_BIN" ]; then
            echo "postinst: sysc-greet already installed, skipping download."
        else
            if command -v wget >/dev/null 2>&1 || command -v curl >/dev/null 2>&1; then
                echo "postinst: Downloading sysc-greet v${SYSC_VERSION}..."

                # Create temp directory for extraction
                SYSC_TMP_DIR="$(mktemp -d /tmp/sysc-greet.XXXXXX)"
                SYSC_DEB_FILE="${SYSC_TMP_DIR}/sysc-greet.deb"

                # Download the deb package
                if command -v curl >/dev/null 2>&1; then
                    curl -sL "$SYSC_DEB_URL" -o "$SYSC_DEB_FILE" 2>/dev/null || true
                else
                    wget -q "$SYSC_DEB_URL" -O "$SYSC_DEB_FILE" 2>/dev/null || true
                fi

                if [ -f "$SYSC_DEB_FILE" ] && [ -s "$SYSC_DEB_FILE" ]; then
                    echo "postinst: Extracting sysc-greet from deb package..."

                    # Extract the deb package
                    dpkg-deb -x "$SYSC_DEB_FILE" "$SYSC_TMP_DIR/extracted" 2>/dev/null || {
                        echo "postinst: WARNING: Failed to extract sysc-greet deb package."
                        rm -rf "$SYSC_TMP_DIR"
                    }

                    # Install the binary (it's in usr/local/bin in the upstream deb)
                    if [ -f "$SYSC_TMP_DIR/extracted/usr/local/bin/sysc-greet" ]; then
                        mkdir -p "$(dirname "$SYSC_BIN")"
                        cp "$SYSC_TMP_DIR/extracted/usr/local/bin/sysc-greet" "$SYSC_BIN"
                        chmod +x "$SYSC_BIN"
                        echo "postinst: sysc-greet binary installed to $SYSC_BIN"

                        # Install data files (ASCII configs, fonts, wallpapers)
                        if [ -d "$SYSC_TMP_DIR/extracted/usr/share/sysc-greet" ]; then
                            mkdir -p "$SYSC_SHARE_DIR"
                            cp -r "$SYSC_TMP_DIR/extracted/usr/share/sysc-greet/"* "$SYSC_SHARE_DIR/"
                            echo "postinst: sysc-greet data files installed to $SYSC_SHARE_DIR"
                        fi

                        # Update greetd config to use sysc-greet instead of wlgreet
                        if [ -f /etc/greetd/config.toml ]; then
                            echo "postinst: Updating greetd configuration to use sysc-greet..."
                            # Backup original config
                            if [ ! -f /etc/greetd/config.toml.bak ]; then
                                cp /etc/greetd/config.toml /etc/greetd/config.toml.bak
                            fi
                            # Replace wlgreet with sysc-greet
                            sed -i 's|command.*wlgreet.*|command = "/usr/local/bin/sysc-greet"|g' /etc/greetd/config.toml
                            echo "postinst: greetd configuration updated."
                        fi

                        # Update sway-config to use sysc-greet instead of wlgreet
                        if [ -f /etc/greetd/sway-config ]; then
                            echo "postinst: Updating greetd sway-config to use sysc-greet..."
                            # Backup original config
                            if [ ! -f /etc/greetd/sway-config.bak ]; then
                                cp /etc/greetd/sway-config /etc/greetd/sway-config.bak
                            fi
                            # Replace wlgreet with sysc-greet and remove "swaymsg exit" for sysc-greet
                            sed -i 's|exec "/usr/sbin/wlgreet; swaymsg exit"|exec "/usr/local/bin/sysc-greet"|g' /etc/greetd/sway-config
                            echo "postinst: greetd sway-config updated."
                        fi

                        echo "postinst: sysc-greet installed successfully."
                        echo "postinst: Note: Please restart greetd to apply changes: systemctl restart greetd"
                    else
                        echo "postinst: WARNING: sysc-greet binary not found in extracted package."
                    fi

                    # Cleanup temp directory
                    rm -rf "$SYSC_TMP_DIR"
                else
                    echo "postinst: WARNING: sysc-greet download failed. Display manager greeter may not work."
                    echo "postinst: You can install it manually from: https://github.com/Nomadcxx/sysc-greet"
                    rm -rf "$SYSC_TMP_DIR"
                fi
            else
                echo "postinst: WARNING: Neither curl nor wget available. Cannot download sysc-greet."
                echo "postinst: Install curl or wget, then reinstall this package."
            fi
        fi

        echo "postinst: ==============================================="
        echo ""

        # Prompt for username to add to sudo group
        echo ""
        echo "postinst: ==============================================="
        echo "postinst: Optional: Add a user to the sudo group"
        echo "postinst: ==============================================="

        SUDO_USERNAME=""

        # Try debconf first (standard Debian way), fallback to read
        if command -v debconf >/dev/null 2>&1; then
            # Try to get username from debconf
            SUDO_USERNAME="$(debconf-show ybarojs-debian 2>/dev/null | grep 'ybarojs-debian/sudo-username' | sed 's/.*: //')" || true
        fi

        # If debconf didn't provide a username, prompt interactively
        if [ -z "$SUDO_USERNAME" ]; then
            # Only prompt if we have a terminal (not during non-interactive installs)
            if [ -t 0 ] && [ -t 1 ]; then
                printf "postinst: Enter username to add to sudo group (leave empty to skip): "
                read -r SUDO_USERNAME || true
            fi
        fi

        # Process the username if provided
        if [ -n "$SUDO_USERNAME" ]; then
            # Validate that the user exists
            if id "$SUDO_USERNAME" >/dev/null 2>&1; then
                echo "postinst: Adding user '$SUDO_USERNAME' to sudo group..."
                if usermod -aG sudo "$SUDO_USERNAME" 2>/dev/null; then
                    echo "postinst: User '$SUDO_USERNAME' successfully added to sudo group."
                    echo "postinst: User must log out and back in for changes to take effect."
                else
                    echo "postinst: WARNING: Failed to add user to sudo group. Run manually: usermod -aG sudo $SUDO_USERNAME"
                fi
            else
                echo "postinst: WARNING: User '$SUDO_USERNAME' does not exist. Skipping sudo group assignment."
                echo "postinst: Create the user first, then run: usermod -aG sudo $SUDO_USERNAME"
            fi
        else
            echo "postinst: No username provided. Skipping sudo group setup."
            echo "postinst: To add a user later, run: usermod -aG sudo <username>"
        fi

        echo "postinst: ==============================================="
        echo ""

        # Add _greetd user to video, render, and input groups for GPU/input access
        echo ""
        echo "postinst: ==============================================="
        echo "postinst: Setting up _greetd user permissions"
        echo "postinst: ==============================================="

        if id "_greetd" >/dev/null 2>&1; then
            for group in video render input; do
                if getent group "$group" >/dev/null 2>&1; then
                    echo "postinst: Adding _greetd to '$group' group..."
                    usermod -aG "$group" "_greetd" 2>/dev/null || true
                else
                    echo "postinst: Group '$group' does not exist, skipping..."
                fi
            done
            echo "postinst: _greetd user permissions configured."
        else
            echo "postinst: WARNING: _greetd user does not exist. Greetd may not be installed."
        fi

        echo "postinst: ==============================================="
        echo ""

        # Enable iwd service for wireless networking
        echo ""
        echo "postinst: ==============================================="
        echo "postinst: Enabling iwd wireless service"
        echo "postinst: ==============================================="

        if command -v systemctl >/dev/null 2>&1; then
            if systemctl enable iwd 2>/dev/null; then
                echo "postinst: iwd service enabled successfully."
                echo "postinst: iwd will start automatically on next boot."
            else
                echo "postinst: WARNING: Could not enable iwd service."
                echo "postinst: You may need to enable it manually: systemctl enable iwd"
            fi
        else
            echo "postinst: WARNING: systemctl not found. Cannot enable iwd service."
        fi

        echo "postinst: ==============================================="
        echo ""

        # Show help message for custom PS1 prompt
        echo ""
        echo "=========================================="
        echo "  Optional: Enable Colored Bash Prompt"
        echo "=========================================="
        echo ""
        echo "To add a beautiful colored prompt to your terminal, run:"
        echo ""
        cat << 'PROMPT_CMD'
  echo 'export PS1="\[\e[1;32m\]\u@\h\[\e[m\]:\[\e[1;34m\]\w\[\e[m\] [\[\e[1;33m\]\$?\[\e[m\]]\$ "' >> ~/.bashrc
PROMPT_CMD
        echo ""
        echo "This will give you:"
        echo "  🟢 Green: username@hostname"
        echo "  🔵 Blue: current directory"
        echo "  🟡 Yellow: exit code of last command"
        echo ""
        echo "Then run: source ~/.bashrc"
        echo ""
        ;;
esac

exit 0
